# dosbox.conf
[program:dosbox]
command=echo "RUN /startgame.sh"
autorestart=false